# crates/ruvector-hailo-cluster/deploy/ruvector-hailo-worker.service
[Unit]
Description=ruvector Hailo embedding worker
Documentation=https://github.com/ruvnet/ruvector
# Keep going until tailscaled is up if we depend on tailnet identity
# (the worker itself doesn't, but operators commonly want it on the
# tailnet for ACL-gated access).
Wants=network-online.target tailscaled.service
After=network-online.target tailscaled.service

[Service]
Type=simple
User=genesis
Group=genesis
WorkingDirectory=/home/genesis
EnvironmentFile=-/etc/ruvector-hailo.env
ExecStart=/usr/local/bin/ruvector-hailo-worker
Restart=on-failure
RestartSec=2

# Sane defaults for an inference daemon — let the kernel kill us before
# we eat the box, and don't let us trample on someone else's CPU sched.
LimitNOFILE=8192
Nice=-5
# Hailo PCIe driver creates /dev/hailo0 owned by root:plugdev — give
# the service supplementary group access without running as root.
SupplementaryGroups=plugdev

# Hardening — don't reach further than embedding inference needs.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/genesis
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
DevicePolicy=closed
DeviceAllow=/dev/hailo0 rw

# Pipe structured logs to journald with rate-limit headroom.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
